# crontab/crontabfile
35 * * * *  /home/pi/rig-monitor/dns_update.sh > /dev/null
* * * * *  /home/pi/rig-monitor/rig-monitor.sh > /dev/null
5-59/10 * * * *  /home/pi/rig-monitor/coin_data.sh > /dev/null
* * * * *  /home/pi/rig-monitor/env-monitor.sh > /dev/null
*/10 * * * *  /home/pi/rig-monitor/pool-monitor.sh > /dev/null
9-59/10 * * * *  /home/pi/rig-monitor/profitability.sh > /dev/null


